# .setup/home-nas-data.automount
[Unit]
Description=Automount odroid media share

[Automount]
Where=/home/nas/data
TimeoutIdleSec=10min

[Install]
WantedBy=multi-user.target
