# .setup/home-nas-nginx.automount
[Unit]
Description=Automount odroid www share

[Automount]
Where=/home/nas/nginx
TimeoutIdleSec=10min

[Install]
WantedBy=multi-user.target
