# setup/nfs/home-nas-data2.automount
[Unit]
Description=Automount odroid media share

[Automount]
Where=/home/nas/data2
TimeoutIdleSec=10min

[Install]
WantedBy=multi-user.target
