# setup/nfs/home-nas-backups.automount
[Unit]
Description=Automount odroid backups share

[Automount]
Where=/home/nas/backups
TimeoutIdleSec=10min

[Install]
WantedBy=multi-user.target
