# setup/nfs/home-nas-media.mount
[Unit]
Description=Mount odroid media share

[Mount]
What=192.168.50.100:/home/nas/media
Where=/home/nas/media
Options=noatime
Type=nfs
TimeoutSec=30
ForceUnmount=true

[Install]
WantedBy=multi-user.target
