# setup/nfs/home-nas-media.automount
[Unit]
Description=Automount odroid media share

[Automount]
Where=/home/nas/media
TimeoutIdleSec=10min

[Install]
WantedBy=multi-user.target
