# setup/nfs/home-nas-data2.mount
[Unit]
Description=Mount odroid media share

[Mount]
What=192.168.50.100:/home/nas/data2
Where=/home/nas/data2
Options=noatime
Type=nfs
TimeoutSec=30
ForceUnmount=true

[Install]
WantedBy=multi-user.target
